# files/wsl2/ubuntu-22/etc/cron.daily/cleanups
#!/bin/bash

# Clean temporary files older than few days
find /tmp /var/tmp -mtime +7 -exec rm -rf {} \;
