# templates/wsl2/ubuntu-22/etc/profile.d/git-performance.sh.j2
#!/bin/bash

# See https://fromkk.com/posts/speed-up-git-speed-in-wsl/
function git() {
  if $(pwd -P | grep -q "^\/mnt\/c\/*"); then
    git.exe "$@"
  else
    command git "$@"
  fi
}
